# .config/systemd/user/change-wallpaper.service
[Unit]
Description=Change wallpaper

[Service]
ExecStart=bash -lc 'change-wallpaper'
