# .config/systemd/user/change-wallpaper.timer
[Unit]
Description=Run change-wallpaper service every 2 minutes after start of this service

[Timer]
OnActiveSec=0
OnUnitActiveSec=2m
