# .config/systemd/user/hyprland.target
[Unit]
Description=Services to be started with hyprland
Wants=change-wallpaper.timer
Wants=cloudflare-tunnel.service
